# files/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/bash -c "[[ -f /home/app/nginx ]] && /home/app/nginx || sleep 3"

[program:app]
command=/home/app/init
